# attense-app/ATTENSE_app/AI/Data/APP-02-CommandInjection.json
{
  "attack_id": "APP-02",
  "attack_name": "Command Injection",
  "owasp": "A03:2021",
  "owasp_source": "https://owasp.org/Top10/A03_2021-Injection/",
  "mitre": "T1059",
  "mitre_source": "https://attack.mitre.org/techniques/T1059/",
  "description": "Execution of arbitrary operating system commands due to improper input validation.",

  "complexity_factors": [
    {
      "factor": "direct_vs_blind",
      "question": "Was the injection blind (no output returned to attacker)?",
      "adds": 1,
      "reason": "Blind command injection leaves fewer traces in web logs because there is no visible output, making it harder for the blue team to confirm exploitation from logs alone.",
      "source": "OWASP A03:2021 — Blind injection variants rated lower detectability than direct output injection"
    },
    {
      "factor": "system_level_command",
      "question": "Did the injected command execute at OS or root level?",
      "adds": 2,
      "reason": "System-level execution means the attacker has full server access. Containment requires host isolation not just application-level blocking.",
      "source": "MITRE T1059 — Command execution at system level maps to high impact, requiring host-level containment per NIST SP 800-61 Section 3.3"
    },
    {
      "factor": "persistence_established",
      "question": "Did the attacker establish persistence via the injected command (e.g. cron job, backdoor)?",
      "adds": 2,
      "reason": "Persistence means containment of the initial vector is not enough — the blue team must also find and remove the persistence mechanism.",
      "source": "MITRE T1053 — Scheduled Task/Job: persistence mechanisms multiply the containment effort required"
    },
    {
      "factor": "obfuscated_command",
      "question": "Was the command obfuscated or encoded to bypass input validation?",
      "adds": 1,
      "reason": "Encoded commands reduce the confidence of signature-based detection tools and require deeper log analysis to confirm.",
      "source": "MITRE T1027 — Obfuscated Files or Information: reduces detection confidence"
    }
  ],

  "ttc_ranges": [
    {
      "complexity_min": 1,
      "complexity_max": 2,
      "acceptable_min_min": 20,
      "acceptable_max_min": 45,
      "max_tolerance_min": 60,
      "source": "NIST SP 800-61 Rev2 — Section 3.3: Low severity incidents should be contained rapidly with minimal escalation"
    },
    {
      "complexity_min": 3,
      "complexity_max": 4,
      "acceptable_min_min": 45,
      "acceptable_max_min": 90,
      "max_tolerance_min": 120,
      "source": "NIST SP 800-61 Rev2 — Section 3.3: Medium severity incidents require analysis before containment"
    },
    {
      "complexity_min": 5,
      "complexity_max": 6,
      "acceptable_min_min": 90,
      "acceptable_max_min": 150,
      "max_tolerance_min": 180,
      "source": "NIST SP 800-61 Rev2 — Section 3.3: High severity incidents involve multiple IR phases before containment is possible"
    }
  ],

  "rules": [
    {
      "id": "CMD-R01",
      "phase": 2,
      "description": "Alert must be investigated within 10 minutes of being raised",
      "type": "major",
      "penalty": -15,
      "auto_fail": false,
      "trigger": "No alert_investigation_started event within 600 seconds of alert_raised",
      "source": "NIST SP 800-61 Rev2 — Section 3.2.4: Initial analysis should begin promptly after detection"
    },
    {
      "id": "CMD-R02",
      "phase": 2,
      "description": "Incident must be confirmed — command execution must be escalated",
      "type": "major",
      "penalty": -15,
      "auto_fail": false,
      "trigger": "No incident_confirmed event in event log",
      "source": "NIST SP 800-61 Rev2 — Section 3.2: Incident declaration is a required step in the IR lifecycle"
    },
    {
      "id": "CMD-R03",
      "phase": 2,
      "description": "Host must be isolated — containment_initiated before max_tolerance",
      "type": "critical",
      "penalty": -25,
      "auto_fail": false,
      "trigger": "No containment_initiated event before max_tolerance_min seconds",
      "source": "NIST SP 800-61 Rev2 — Section 3.3: For system-level compromise, host isolation is the primary containment strategy"
    },
    {
      "id": "CMD-R04",
      "phase": 2,
      "description": "Containment must succeed — host isolated or command execution blocked",
      "type": "critical",
      "penalty": -30,
      "auto_fail": true,
      "trigger": "No containment_succeeded event at all",
      "source": "NIST SP 800-61 Rev2 — Section 3.3: Containment must be confirmed before incident can be closed"
    }
  ],

  "scoring": {
    "base_score": 100,
    "verdict_thresholds": {
      "excellent":    { "min": 90 },
      "acceptable":   { "min": 70 },
      "needs_review": { "min": 50 },
      "failed":       { "min": 0  }
    }
  },

  "scenarios": [
    {
      "scenario_id": "CMD-S1",
      "label": "Direct command injection — good response",
      "description": "Blue team detects abnormal process execution, confirms incident, and isolates the host within acceptable time.",
      "complexity_factors_active": ["system_level_command"],
      "complexity": 2,
      "expected_score": 85,
      "expected_verdict": "acceptable",
      "logs": [
        { "t_offset_sec": 0,    "actor_type": "system",    "actor_id": "wazuh",     "event_type": "alert_raised",                "target_type": "host",    "outcome": "success", "detail": "Wazuh rule 31120 fired — OS command separator detected in web request" },
        { "t_offset_sec": 90,   "actor_type": "blue_team", "actor_id": "analyst-1", "event_type": "alert_investigation_started", "target_type": "alert",   "outcome": "success", "detail": "Analyst-1 opened alert, began process log review" },
        { "t_offset_sec": 480,  "actor_type": "blue_team", "actor_id": "analyst-2", "event_type": "incident_confirmed",          "target_type": "host",    "outcome": "success", "detail": "Analyst-2 confirmed abnormal child process spawned by web server" },
        { "t_offset_sec": 2400, "actor_type": "blue_team", "actor_id": "analyst-3", "event_type": "containment_initiated",       "target_type": "host",    "outcome": "success", "detail": "Analyst-3 initiated host network isolation" },
        { "t_offset_sec": 2700, "actor_type": "blue_team", "actor_id": "analyst-3", "event_type": "containment_succeeded",       "target_type": "host",    "outcome": "success", "detail": "Host isolated from network, command execution blocked" }
      ]
    },
    {
      "scenario_id": "CMD-S2",
      "label": "Blind command injection with persistence — missed containment",
      "description": "Blue team detects alert but fails to contain before tolerance window. Persistence mechanism not identified. Auto-fail triggered.",
      "complexity_factors_active": ["direct_vs_blind", "persistence_established", "obfuscated_command"],
      "complexity": 4,
      "expected_score": 35,
      "expected_verdict": "failed",
      "logs": [
        { "t_offset_sec": 0,    "actor_type": "system",    "actor_id": "wazuh",     "event_type": "alert_raised",                "target_type": "host",  "outcome": "success", "detail": "Wazuh rule 31121 fired — encoded command pattern detected" },
        { "t_offset_sec": 300,  "actor_type": "blue_team", "actor_id": "analyst-1", "event_type": "alert_investigation_started", "target_type": "alert", "outcome": "success", "detail": "Analyst-1 opened alert" },
        { "t_offset_sec": 1200, "actor_type": "blue_team", "actor_id": "analyst-1", "event_type": "incident_confirmed",          "target_type": "host",  "outcome": "success", "detail": "Analyst-1 confirmed suspicious process" }
      ]
    }
  ]
}
