# blueteam/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

; ── BlueTeam FastAPI (the only public-facing process) ────────────────────────
[program:blueteam-api]
command=uvicorn main:app --host 0.0.0.0 --port 8010
directory=/app
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/blueteam-api.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/blueteam-api.err
stderr_logfile_maxbytes=5MB

; ── Cassandra (TheHive storage) ───────────────────────────────────────────────
[program:cassandra]
command=/usr/sbin/cassandra -f
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/cassandra.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/cassandra.err
stderr_logfile_maxbytes=10MB

; ── Elasticsearch (TheHive index) ────────────────────────────────────────────
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
environment=ES_JAVA_OPTS="-Xms512m -Xmx512m",discovery.type="single-node",xpack.security.enabled="false"
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/elasticsearch.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/elasticsearch.err
stderr_logfile_maxbytes=10MB

; ── TheHive (analyst UI + case management) ───────────────────────────────────
[program:thehive]
command=/opt/thehive/bin/thehive -Dconfig.file=/etc/thehive/application.conf
autostart=true
autorestart=true
priority=30
; TheHive must wait for Cassandra and Elastic — start.sh handles the wait
startsecs=30
stdout_logfile=/var/log/thehive.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/thehive.err
stderr_logfile_maxbytes=10MB
