# attense-app/ATTENSE_app/blueteam/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

; ── BlueTeam FastAPI (the only public-facing process) ────────────────────────
; TheHive, Cassandra and Elasticsearch are NOT installed in this container.
; The HiveClient stub handles unreachable-Hive gracefully (returns {} on error).
; Only the FastAPI service runs here.
[program:blueteam-api]
command=uvicorn main:app --host 0.0.0.0 --port 8010
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/blueteam-api.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/blueteam-api.err
stderr_logfile_maxbytes=5MB
